# systemd/omarchy-theme-switcher.timer
[Unit]
Description=Omarchy Theme Switcher — boot anchor
Documentation=https://github.com/lukas/omarchyThemeSwitcher
After=graphical-session.target

[Timer]
OnBootSec=30s
Persistent=true
AccuracySec=10s

[Install]
WantedBy=timers.target
